# clamd.scanonaccess.init
#!/bin/bash
#
# chkconfig: - 78 32
# description: The clamd server running for amavisd

CLAMD_SERVICE=scanonaccess
. /usr/share/clamav/clamd-wrapper
