# clamd@scanonaccess.service
.include /usr/lib/systemd/system/clamd@.service

[Unit]
Description = ScanOnAccess clamav scanner daemon

[Install]
WantedBy = multi-user.target
